# systemd/alacritty_master.timer
[Unit]
Description=Run Alacritty master benchmark daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
