# systemd/alacritty_master.service
[Unit]
Description=Benchmark Alacritty master performance
After=network.target

[Service]
Type=oneshot
ExecStart=/home/perfbot/termbenchbot/alacritty_master.sh
WorkingDirectory=/home/perfbot/termbenchbot/

[Install]
WantedBy=multi-user.target
