# systemd/plot_master.timer
[Unit]
Description=Plot Alacritty master performance daily

[Timer]
OnCalendar=*-*-* 1:00:00
Persistent=true

[Install]
WantedBy=timers.target
