# systemd/plot_master.timer
[Unit]
Description=Plot Alacritty master performance daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
